# utils/flight/dasng.service
[Unit]
Link Engineering Data Acquisition System ng Service

[Service]
ExecStart=/usr/local/sbin/flight.sh

[Install]
WantedBy=multi-user.target
